# WebAssignment/logTemp.cron
0,15,30,45 * * * *  /usr/bin/bin/python /home/pi/ELSpring2018/WebAssignment/ReadTemp.py

