# Final_Project/Get_data.cron
*/5 * * * * python /home/pi/ELSpring2018/Final_Project/dht_read.py